# scripts/macro-zero-run.service
[Unit]
Description=Macro-Zero Service
Requires=macro-zero-startup.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/bin/macro-zero/src/macrozero.py
User=pi
StartLimitIntervalSec=10
Restart=on-failure

[Install]
WantedBy=multi-user.target
